# projects/character-lcd-display-linux/projMBLCD03/Linux/mblcd03/commandline/etc-gentoo/init.d/mblcd
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

PIDFILE=/var/run/mblcd.pid
DAEMON=/usr/local/bin/mblcd03d

depend() {
	use usb
}

start() {
	ebegin "Starting mblcd"
	start-stop-daemon --start --pidfile ${PIDFILE} --exec ${DAEMON} -- ${MBLCD_OPTS}
	eend $? "Failed to start mblcd."
}

stop() {
	ebegin "Stopping mblcd"
	start-stop-daemon --stop --quiet --oknodo --pidfile ${PIDFILE} --exec ${DAEMON}
	eend $?
}
